# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:zookeeper]
command=/usr/local/zookeeper/bin/zkServer.sh start-foreground
user=daemon
priority=0

[program:mysql]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe
user=mysql
priority=0

[program:druid-realtime]
user=druid
command=java
  -server
  -Xmx1g
  -Ddruid.service=realtime
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8

  -Ddruid.host=%(ENV_HOSTIP)s  
  -Ddruid.port=8083

  -Ddruid.zk.service.host=localhost  

  -Ddruid.extensions.coordinates=["io.druid.extensions:druid-examples:0.6.171","io.druid.extensions:druid-kafka-seven:0.6.171"]

  -Ddruid.realtime.specFile=/usr/local/druid/spec_file.spec
  -Ddruid.publish.type=noop

  -Ddruid.processing.buffer.sizeBytes=100000000
  -Ddruid.processing.numThreads=1

  -cp /usr/local/druid/lib/*  
  io.druid.cli.Main server realtime
autorestart=true
redirect_stderr=true
priority=100

[program:druid-coordinator]
user=druid
command=java
  -server
  -Xmx1g
  -Ddruid.service=coordinator
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Ddruid.zk.service.host=localhost
  -Ddruid.host=%(ENV_HOSTIP)s  

  -Ddruid.db.connector.connectURI=jdbc\:mysql\://localhost\:3306/druid
  -Ddruid.db.connector.user=druid
  -Ddruid.db.connector.password=diurd

  -Ddruid.coordinator.startDelay=PT5S
  -cp /usr/local/druid/lib/*
  io.druid.cli.Main server coordinator
autorestart=true
redirect_stderr=true
priority=100

[program:druid-historical]
user=druid
command=java
  -server
  -Xmx1g
  -Ddruid.service=historical
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Ddruid.zk.service.host=localhost
  -Ddruid.host=%(ENV_HOSTIP)s    
  -Ddruid.storage.type=c*
  -Ddruid.storage.host=localhost:9160
  -Ddruid.storage.keyspace=druid  
  -Ddruid.computation.buffer.size=67108864
  -Ddruid.segmentCache.locations="[{\"path\":\"/var/tmp/druid/indexCache\",\"maxSize\":5000000000}]"
  -Ddruid.server.maxSize=5000000000
  -cp /usr/local/druid/lib/*
  io.druid.cli.Main server historical
autorestart=true
redirect_stderr=true
priority=100

[program:druid-broker]
user=druid
command=java
  -server
  -Xmx1g
  -Ddruid.service=broker
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Ddruid.zk.service.host=localhost
  -Ddruid.host=%(ENV_HOSTIP)s
  -Ddruid.computation.buffer.size=67108864
  -Ddruid.broker.cache.sizeInBytes=33554432
  -cp /usr/local/druid/lib/*
  io.druid.cli.Main server broker
autorestart=true
redirect_stderr=true
priority=100
